# .cron
# buy-crypto: Runs every day at 9:00 AM
0 9 * * * /home/ec2-user/.bun/bin/bun /home/ec2-user/app/upbit-trade/batch/runner/runBuyCrypto.ts

# sell-crypto: Runs every second
* * * * * /home/ec2-user/.bun/bin/bun /home/ec2-user/app/upbit-trade/batch/runner/runSellCrypto.ts

# swap: Runs every day at 10:00 PM
0 22 * * * /home/ec2-user/.bun/bin/bun /home/ec2-user/app/upbit-trade/batch/runner/runSwap.ts

# health-check: Runs every day at 9:00 AM
0 9 * * * /home/ec2-user/.bun/bin/bun /home/ec2-user/app/upbit-trade/batch/runner/runHealthCheck.ts